# systemd/soundboard.socket
[Unit]
Description=Soundboard varlink socket

[Socket]
ListenStream=%t/soundboard.varlink
Accept=no

[Install]
WantedBy=sockets.target
